# files/supervisor/wildfly.ini
[program:wildfly]
user = wildfly
autostart = true
command = /opt/wildfly/bin/standalone.sh -b 0.0.0.0 -bmanagement=0.0.0.0
redirect_stderr = true
